# docs/general-notes/scripts/tailscale-restart.service
[Unit]
Description=Restart Tailscale after startup
After=tailscaled.service network-online.target
Wants=network-online.target
Requires=tailscaled.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/root/tailscale-restart.sh
User=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target